# .service
[Unit]
Description=Gunicorn instance to serve joy-app
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/joy-app
Environment="PATH=/home/ubuntu/joy-app/env/bin"
ExecStart=/home/ubuntu/joy-app/env/bin/gunicorn --workers 3 --bind unix:joy-app.sock -m 007 wsgi:app

[Install]
WantedBy=multi-user.target